# infrastructure/systemd/kraken-bot.service
[Unit]
Description=AlgonikHQ Kraken DCA Bot
Documentation=https://github.com/AlgonikHQ/trading-stack
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kraken_bot

EnvironmentFile=/root/kraken_bot/.env

ExecStart=/usr/bin/python3 /root/kraken_bot/bot.py

Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kraken-bot

TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
