# infrastructure/systemd/oanda-bot.service
[Unit]
Description=AlgonikHQ OANDA Trading Bot
Documentation=https://github.com/AlgonikHQ/trading-stack/tree/main/oanda-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
EnvironmentFile=/root/oanda_bot.env
ExecStart=/usr/bin/python3 /root/oanda_trading_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
