# infrastructure/systemd/scalp-sniper.service
[Unit]
Description=AlgonikHQ OSC Bot
Documentation=https://github.com/AlgonikHQ/trading-stack/tree/main/solana-sniper/osc-scalper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/solana_sniper
EnvironmentFile=/root/solana_sniper/.env
ExecStart=/usr/bin/python3 -u /root/solana_sniper/scalp_scanner.py
Restart=on-failure
RestartSec=30
KillMode=mixed
StandardOutput=null
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
