# infrastructure/systemd/solana-sniper.service
[Unit]
Description=AlgonikHQ Solana Sniper Bot
Documentation=https://github.com/AlgonikHQ/trading-stack/tree/main/solana-sniper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/solana_sniper
EnvironmentFile=/root/solana_sniper/.env
ExecStart=/usr/bin/python3 -u /root/solana_sniper/scanner.py
ExecStopPost=/bin/bash -c 'curl -s -X POST "https://api.telegram.org/bot${TELEGRAM_BOT_TOKEN}/sendMessage" -d chat_id="${TELEGRAM_CHAT_ID}" -d text="🔴 Solana Sniper STOPPED"'
Restart=always
RestartSec=30
KillMode=mixed
StandardOutput=append:/var/log/solana_sniper.log
StandardError=append:/var/log/solana_sniper.log

[Install]
WantedBy=multi-user.target
